# ansible/04_k3s/roles/storage/templates/syncoid.timer
[Unit]
Description=Runs syncoid

[Timer]
OnCalendar=*-*-* *:10:00
Unit=syncoid.service

[Install]
WantedBy=timers.target
